# coffeeshop-security-lab/coffeeshop-security-lab/docker/jokopi/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
autostart=true
autorestart=true
priority=10

[program:php-fpm]
command=/usr/sbin/php-fpm%(ENV_PHP_VERSION)s -F
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm-error.log

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log

[program:wazuh-agent]
command=/var/ossec/bin/wazuh-control start
autostart=true
autorestart=true
priority=40
startsecs=10
startretries=5
stdout_logfile=/var/log/wazuh-agent.log
stderr_logfile=/var/log/wazuh-agent-error.log
