# docker/jokopi/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
autostart=true
autorestart=true
priority=10

[program:link-php-fpm-sock]
# Creates /run/php/php-fpm.sock as a symlink to whichever versioned socket
# (e.g. /run/php/php8.1-fpm.sock) PHP-FPM ends up creating. nginx.conf
# references the un-versioned name; without this shim the FastCGI proxy
# fails and every /api/*.php request returns 404. See the long comment
# in the Dockerfile near the symlink helper definition.
#
# priority=15 puts this AFTER rsyslog (10) but BEFORE php-fpm (20).
# startsecs=0 because the helper exits as soon as the socket appears
# (or after a 10s timeout); supervisord doesn't need to keep it alive.
command=/usr/local/bin/link-php-fpm-sock
autostart=true
autorestart=false
startsecs=0
exitcodes=0
priority=15
stdout_logfile=/var/log/link-php-fpm-sock.log
stderr_logfile=/var/log/link-php-fpm-sock-error.log

[program:php-fpm]
# The original config used `/usr/sbin/php-fpm%(ENV_PHP_VERSION)s -F` to
# interpolate the PHP version at runtime, but nothing ever sets PHP_VERSION
# in the container env, so supervisord fails to parse this section and the
# container restart-loops with exit code 2:
#   "Format string '/usr/sbin/php-fpm%(ENV_PHP_VERSION)s -F' for
#    'program:php-fpm.command' contains names ('ENV_PHP_VERSION') which
#    cannot be expanded."
# The Dockerfile now creates a stable `/usr/sbin/php-fpm` symlink that
# points at whichever versioned binary apt installed, so we can reference
# it directly without any runtime interpolation.
command=/usr/sbin/php-fpm -F
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm-error.log

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log

[program:wazuh-agent]
# `wazuh-control start` spawns the agent's daemons (ossec-agentd,
# wazuh-modulesd, ossec-syscheckd, ossec-execd, ossec-logcollector) as
# detached background processes and then exits with code 0. Under the
# previous config (autorestart=true, startsecs=10) supervisord treated
# the clean <10 s exit as an "unexpected termination" and kept relaunching
# `wazuh-control start`, which printed "Process is already running" for
# each daemon and looked, in the logs, exactly like a crashloop.
#
# Two changes fix it:
#   1. Wrap the start in `sh -c "wazuh-control start && tail -F ossec.log"`
#      so the supervisord-managed process is a long-lived tail. That
#      gives supervisord a real foreground PID to watch AND surfaces the
#      agent log in `docker compose logs jokopi-app`.
#   2. Drop startsecs back to 0 (the start command's own time is what
#      matters; the tail is instant) and keep autorestart=true so that
#      if the tail dies for any reason supervisord reattaches.
command=/bin/sh -c "/var/ossec/bin/wazuh-control start && exec tail -F /var/ossec/logs/ossec.log"
autostart=true
autorestart=true
startsecs=0
priority=40
stdout_logfile=/var/log/wazuh-agent.log
stderr_logfile=/var/log/wazuh-agent-error.log
